# misc/webcam.crontab
# Purge webcam images more than 5 minutes old from the staging directory.
*/5 * * * * find /tmp/webcam -type f -cmin +5 -delete
